# deploy/vk-parser.service
[Unit]
Description=VK Community Post Parser
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vkparser
Group=vkparser
WorkingDirectory=/opt/vk-community-parser
EnvironmentFile=/opt/vk-community-parser/.env
ExecStart=/opt/vk-community-parser/.venv/bin/python /opt/vk-community-parser/main.py
Restart=always
RestartSec=15
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
